# server/crontab
* * * * * source /watchman/server/run_all_jobs.sh